# etc/radiusd.conf
[unix_http_server]
file=/tmp/toughradius.sock

[inet_http_server]
port=127.0.0.1:19001
username=ctlman
password=ctlroot

[supervisord]
nodaemon=false
logfile=/etc/toughradius/toughradius.log
logfile_maxbytes=1MB
logfile_backups=8
loglevel=debug
pidfile=/etc/toughradius/toughradius.pid
;user=root

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/toughradius.sock

[program:api]
command=gtrad apiserv -c /etc/toughradius/radiusd.json
dictionary=/etc/toughradius
startretries = 5
autorestart = true
redirect_stderr=true
stdout_logfile=/etc/toughradius/gtrad-api.log


[program:auth]
command=gtrad auth -c /etc/toughradius/radiusd.json
dictionary=/etc/toughradius
startretries = 5
autorestart = true
redirect_stderr=true
stdout_logfile=/etc/toughradius/gtrad-auth.log

[program:acct]
command=gtrad acct -c /etc/toughradius/radiusd.json
dictionary=/etc/toughradius
startretries = 5
autorestart = true
redirect_stderr=true
stdout_logfile=/etc/toughradius/gtrad-acct.log

