# etc/supervisord_test.conf
[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
nodaemon=true
logfile=/tmp/supervisord.log
logfile_maxbytes=8MB
logfile_backups=3
loglevel=info
pidfile=/tmp/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock



[program:manage]
command=pypy toughctl --manage -c etc/toughradius.json
redirect_stderr=true
stdout_logfile=/tmp/toughradius.log


[program:auth]
command=pypy toughctl --auth -c etc/toughradius.json
redirect_stderr=true
stdout_logfile=/tmp/toughradius.log

[program:acct]
command=pypy toughctl --acct -c etc/toughradius.json
redirect_stderr=true
stdout_logfile=/tmp/toughradius.log

[program:worker]
command=pypy toughctl --worker -c etc/toughradius.json
process_name = %(program_name)s%(process_num)d
numprocs=4
redirect_stderr=true
stdout_logfile=/tmp/toughradius.log

[program:task]
command=pypy toughctl --task -c etc/toughradius.json
redirect_stderr=true
stdout_logfile=/tmp/toughradius.log