# etc/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock

[inet_http_server]
port=127.0.0.1:9111

[supervisord]
logfile=/var/log/supervisord.log
logfile_maxbytes=8MB
logfile_backups=3
loglevel=info
pidfile=/tmp/supervisord.pid
nodaemon=false
minfds=1024
minprocs=200

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock
serverurl=http://127.0.0.1:9111

[program:admin]
command=pypy /opt/toughradius/toughctl --admin
process_name = %(program_name)s%(process_num)d
numprocs = 1
directory=/opt/toughradius
autostart=true
autorestart=true
user=root
redirect_stderr=true
stdout_logfile=/var/toughradius/log/admin.log


