# install/backupdb.cron
30 1 * * * /usr/bin/toughrad backupdb > /dev/null
